# deploy/com.immie.yin.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.immie.yin</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.14/bin/python3</string>
        <string>/Users/immie/Projects/Yin v2/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/immie/Projects/Yin v2</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/immie/Projects/Yin v2/runtime_state/yin-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/immie/Projects/Yin v2/runtime_state/yin-launchd.log</string>
</dict>
</plist>
